# com.claude.ai-research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.ai-research</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zalen/.codex/bin/ai-research-dashboard-runner.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/zalen/Documents/Zalen/AI Wif Brain Projects/ai-research/logs/launchagent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zalen/Documents/Zalen/AI Wif Brain Projects/ai-research/logs/launchagent.log</string>
</dict>
</plist>
