# deploy/ai-dc-etl.service
[Unit]
Description=AI & DC Dashboard deterministic ETL refresh
Documentation=https://github.com/cutout-z/ai-dc-dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ai-dc
Group=ai-dc
WorkingDirectory=/opt/ai-dc-dashboard
EnvironmentFile=-/etc/ai-dc-dashboard/etl.env
ExecStart=/opt/ai-dc-dashboard/deploy/run-vps-etl.sh
TimeoutStartSec=2h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
StandardOutput=journal
StandardError=journal

